# config/systemd/lecture-tools.service
[Unit]
Description=Lecture Tools legacy web server (retired)
After=network.target

[Service]
Type=oneshot
# The legacy FastAPI service has been removed. This unit intentionally does nothing
# and serves purely as documentation for historical deployments.
ExecStart=/usr/bin/false

[Install]
WantedBy=multi-user.target
